# spark/supervisord.worker.conf
[supervisord]
nodaemon=true

[program:spark]
command=/bin/bash -c "/home/spark/sbin/start-slave.sh spark://`ping master1.spark_master.dev.docker -c 1 | grep from | cut -d' ' -f 4`7077"

