# spark/supervisord.master.conf
[supervisord]
nodaemon=true

[program:master]
command=/bin/bash -c "/home/spark/sbin/start-master.sh"

